# packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        # Remove systemd service if installed
        if command -v systemctl >/dev/null 2>&1; then
            if systemctl is-enabled dck-bootstrap 2>/dev/null | grep -q "enabled"; then
                systemctl disable dck-bootstrap 2>/dev/null || true
            fi
            if systemctl is-active dck-bootstrap 2>/dev/null | grep -q "active"; then
                systemctl stop dck-bootstrap 2>/dev/null || true
            fi
            systemctl daemon-reload 2>/dev/null || true
        fi
        ;;
esac

exit 0
